# .config/systemd/user/dotfiles-update-check.timer
[Unit]
Description=Periodic dotfiles update check

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
Persistent=true
Unit=dotfiles-update-check.service

[Install]
WantedBy=timers.target
